# files/etc/uci-defaults/10-test_podman
echo run hello-world >> /dev/kmsg
podman run -it hello-world >> /dev/kmsg
