# files/etc/uci-defaults/99-shutdown
# this is only for testing

echo Force shurdonw >> /dev/kmsg
reboot -fn
