# files/etc/uci-defaults/99-shutdown
# this is only for testing

df -h >> /dev/kmsg

podman --version >> /dev/kmsg

podman info >> /dev/kmsg 2>> /dev/kmsg

echo Force shutdown >> /dev/kmsg
reboot -fn
